# public/atlas/reviews/47558c80121155dd.json
{"one_line_takeaway":"TTPDetect combines hybrid (dense + LLM neural) retrieval, on-demand call-graph context exploration, and pre-synthesized TTP-specific reasoning guidelines to map stripped malware binaries to MITRE ATT&CK TTPs, reaching 93.25%/93.81% precision/recall on function-level attribution and 85.7% report coverage plus 90.0% precision on newly discovered TTPs for real-world binaries with expert reports.","key_finding":"On function-level recognition (RQ1), TTPDetect achieves average precision 93.25%, recall 93.81%, F1 92.61% across 14 TTPs, versus Baseline I (72.28% / 68.45% / 68.59%) and Baseline II (82.87% / 75.03% / 76.34%) -- outperforming baselines by 10.38 points precision and 18.78 points recall (paper's own stated deltas vs. the stronger Baseline II). TTPDetect reaches 100% precision and recall on T1059, T1112, T1055, T1562, and T1105; its largest baseline margin is on T1036 (Masquerading), F1 86.96% vs. 40.00% (Baseline I) and 36.36% (Baseline II). On binary-level recognition (RQ2, Dataset II), TTPDetect predicts 95 TTPs total with 83 confirmed true positives (87.37% precision); against the two binaries with expert reports it recovers 6 of 7 documented TTPs (85.7% report coverage) and additionally surfaces 24 previously unreported TTPs, of which 21 (90.0%) are validated true positives by human review, averaging 10.5 net-new TTPs per malware sample. Retrieval-stage ablation (Fig. 7) shows neural retrieval cuts candidate function-TTP pairs by 43.65% relative to dense retrieval alone; adding score-based filtering yields 65.66% reduction; the full hybrid pipeline yields an 85.89% reduction relative to dense retrieval alone. Component ablation (RQ4): removing the Context Explorer drops recall sharply (93.81% -> 74.30%) with a smaller precision drop (93.25% -> 85.43%), while removing the Reasoning Guideline drops precision sharply (93.25% -> 77.83%) with a smaller recall drop (93.81% -> 86.63%) -- indicating the two components play complementary, largely orthogonal roles. Cross-model generalization: gpt-4o-2024-08-06 achieves 76.40% precision / 74.5% recall and gpt-4_1-2025-04-14 achieves 89.50% precision / 77.90% recall in the same pipeline, both below the Claude 3.7 Sonnet configuration.","core_contribution":"The paper presents TTPDetect, described as the first LLM agent for recognizing MITRE ATT&CK Tactics, Techniques, and Procedures (TTPs) directly from stripped malware binaries at both the function and binary level. The pipeline has three stages: (1) an LLM-based renaming agent that assigns human-readable names/summaries to decompiled functions in bottom-up call-graph order; (2) a hybrid retrieval stage that narrows the space of candidate (function, TTP) pairs by combining embedding-based dense retrieval (OpenAI text-embedding-3-large + FAISS) with LLM-based 'neural retrieval' that generates an over-inclusive candidate TTP list per function plus a confidence score for pruning; and (3) a function-level 'analyzing agent' composed of a Context Explorer (an agentic tool-use loop with retrieve_function/retrieve_caller primitives that incrementally pulls in caller/callee code on demand) and a TTP-Specific Reasoning Guideline (a structured, per-TTP checklist of core components, positive/negative indicators, and differentiation criteria pre-synthesized offline from the ATT&CK knowledge base) that together produce a binary presence/absence decision for each candidate function-TTP pair. The authors also construct two new datasets: Dataset I (function-level TTP labels derived from 18 leaked C/C++ malware source projects, compiled and stripped, spanning 12 families/8 categories/42 TTPs/210 function-TTP pairs across 122 functions) and Dataset II (8 real-world binaries from MalwareBazaar/VirusShare, 2 with expert-written analysis reports).","threat_model":"The paper does not propose a new attack; it is a defensive/analysis tool. Its operating assumption is the standard malware-analyst setting: an analyst possesses a stripped malware binary (symbols and debug info removed, often with function counts inflated 10-40x relative to source due to compiler-generated helpers and static linking) and needs to determine which ATT&CK TTPs the binary implements, without relying on dynamic execution (which the paper argues is evadable via logic bombs/environment checks and heavyweight to scale) or on pre-existing CTI reports/source code (which are often unavailable). The system uses IDA Pro 7.6 for decompilation and Claude 3.7 Sonnet (via AutoGen v0.6.1) as the underlying LLM, with MITRE ATT&CK v16.1 as the knowledge base.","limitations":"Authors' own stated threats to validity: (1) Annotation correctness -- TTP ground truth is manually assigned by the authors from source code (mitigated by two independent annotators with evidence citations and disagreement resolution, but still author-generated labels, not third-party ground truth); (2) dataset diversity -- Dataset I covers only 18 leaked C/C++ projects and Dataset II only 8 binaries, with only 2 possessing expert reports for report-coverage evaluation, which the authors acknowledge 'may not fully represent the diversity of real-world malware'; (3) potential data leakage from LLM pretraining -- partially mitigated by including two malware families (Umbrellastand, DamascenedPeacock) that first appeared after the evaluated LLMs' training cutoffs, but this yields a very small sample (2 binaries) for the leakage-controlled report-coverage metric; (4) cross-model generalization results show substantial performance variation (precision ranging 76.40%-93.25%) depending on the underlying LLM, meaning reported headline numbers are somewhat model-specific (tied to Claude 3.7 Sonnet); (5) the paper itself documents a concrete failure mode where TTPDetect misclassifies scheduled-task-based persistence as T1547 (Boot or Logon Autostart Execution) instead of the correct T1053 (Scheduled Task/Job) because the reasoning guideline's differentiation criteria do not capture that subtle distinction.","facts":{"datasets_used":["Dataset I: 18 leaked malware C/C++ source projects, compiled and stripped, 12 malware families across 8 threat categories and 2 platforms (Linux, Windows), 42 distinct TTPs, 210 function-TTP pairs over 122 functions","RQ1 evaluation subset of Dataset I: 14 TTPs (each appearing >=5 times), 140 function-TTP pairs across 93 functions, balanced with 59 no-TTP functions and 29 out-of-set-TTP functions for 181 total functions and 2,534 function-TTP classification instances","Dataset II: 8 real-world malware binaries from MalwareBazaar and VirusShare, spanning 203 ATT&CK TTPs considered; only 2 binaries (Umbrellastand, DamascenedPeacock) have expert-written analysis reports for report-coverage evaluation","MITRE ATT&CK v16.1 used as the TTP knowledge base for reasoning guideline synthesis"],"quantitative_results":["Function-level (RQ1): TTPDetect avg precision 93.25%, recall 93.81%, F1 92.61% across 14 TTPs","TTPDetect achieves 100% precision and recall on T1059, T1112, T1055, T1562, T1105","T1036 (Masquerading) F1: TTPDetect 86.96% vs Baseline I 40.00% vs Baseline II 36.36%","T1547 (Boot or Logon Autostart Execution) is TTPDetect's weakest TTP: 62.50% precision / 100% recall / 76.92% F1, due to a documented confusion with T1053 (Scheduled Task/Job)","Binary-level (RQ2, Dataset II): 95 TTPs predicted total, 83 confirmed true positive, overall precision 87.37%","Report coverage: 6 of 7 expert-report-documented TTPs recovered = 85.7% coverage (3/3=100% for Umbrellastand, 3/4=75.0% for DamascenedPeacock)","Newly discovered TTPs beyond expert reports: 24 discovered, 21 confirmed true positive by human validation = 90.0% precision on discoveries, averaging 10.5 previously-unreported TTPs per malware","Candidate-pair reduction (Fig. 7): neural retrieval reduces candidates by 43.65% vs dense retrieval alone; neural w/ score filtering by 65.66%; full TTPDetect pipeline by 85.89%","Ablation - w/o Context Explorer: precision 85.43% (from 93.25%), recall 74.30% (from 93.81%), F1 77.66% (from 92.61%) -- recall drops most","Ablation - w/o Reasoning Guideline: precision 77.83% (from 93.25%), recall 86.63% (from 93.81%), F1 80.52% (from 92.61%) -- precision drops most","Cross-model generalization: gpt-4o-2024-08-06 = 76.40% precision / 74.5% recall; gpt-4_1-2025-04-14 = 89.50% precision / 77.90% recall (both below the Claude 3.7 Sonnet main-result configuration)"],"baselines_compared":["Baseline I: direct/basic LLM prompting on raw stripped decompiled functions (no renaming/summarization), derived from prior CTI-report and interpreted-package TTP-mapping prompting approaches -- avg precision 72.28%, recall 68.45%, F1 68.59%","Baseline II: same direct prompting approach as Baseline I but applied to functions augmented with LLM-generated names and summaries -- avg precision 82.87%, recall 75.03%, F1 76.34%","Retrieval-only ablations for candidate-pair volume (RQ2/Fig.7): dense retrieval alone, neural retrieval alone, and neural retrieval with confidence-score filtering, each compared against the full TTPDetect hybrid pipeline"]},"badge":{"generator_model":"claude-sonnet-5","verifier_model":"claude-sonnet-5","verified_at":"2026-07-19T13:22:25Z","verdict":"clean"}}